# conf/start/zentral/probes/osx-attacks.json
{
  "name": "osx-attacks",
  "osquery": {
    "schedule": [
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'apple.launchd.plist' or name = 'com.apple.launchport.plist';",
        "description": "(http://securelist.com/blog/research/57331/the-icefog-apt-a-tale-of-cloak-and-three-daggers/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.FolderActionsxl.plist';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_mackontrol_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.vsearch.agent.plist' or name = 'com.vsearch.daemon.plist' or name = 'com.vsearch.helper.plist' or name = 'Jack.plist';",
        "description": "(http://www.thesafemac.com/arg-downlite/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from startup_items where path like '%iWorkServices%';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_iworkserv_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where path like '%UserEvent.System.plist';",
        "description": "(http://www.intego.com/mac-security-blog/new-mac-trojan-discovered-related-to-syria/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.BT.BPK.plist';",
        "description": "(http://www.blazingtools.com/mac_keylogger.html)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.genieo.completer.download.plist' or name = 'com.genieo.completer.update.plist' or name = 'com.genieo.completer.ltvbit.plist' or name = 'com.installer.completer.download.plist' or name = 'com.installer.completer.update.plist' or name = 'com.installer.completer.ltvbit.plist' or name = 'com.genieoinnovation.macextension.plist' or name = 'com.genieoinnovation.macextension.client.plist' or name = 'com.genieo.engine.plist';",
        "description": "(http://www.thesafemac.com/arg-genieo/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.webhelper.plist' or name = 'com.webtools.update.agent.plist' or name = 'com.webtools.uninstaller.plist';",
        "description": "(http://www.thesafemac.com/arg-buca-apps/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from apps where bundle_identifier = 'com.ht.RCSMac' or bundle_package_type like 'OSAX';",
        "description": "Detect RAT used by Hacking Team",
        "interval": "86400"
      },
      {
        "value": "Potential information stealing/leakage.",
        "query": "select label, description, authorizations, path, count(path) as collisions from keychain_acls where label != '' and path != '' group by label having collisions > 1;",
        "description": "Detect OS X keychain items with ACLs permitting multiple application access.",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.whitesmoke.uploader.plist';",
        "description": "(http://www.thesafemac.com/osxfkcodec-a-in-action/ )",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.DockActions.plist' or name like '%www. google.com.tstart.plist%';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_olyx_c.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.pwned.plist' or name = 'com.openbundle.plist' or name = 'com.adobe.reader.plist';",
        "description": "(https://www.f-secure.com/v-descs/inqtana_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where path like '%com.conduit.loader.agent.plist' or name = 'com.conduit.loader.agent.plist' or path like '%com.perion.searchprotectd.plist' or name = 'com.perion.searchprotectd.plist';",
        "description": "(http://www.thesafemac.com/arg-conduit/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.PubSabAgent.plist';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_sabpab_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where label = 'com.ht.RCSMac' or label like 'com.yourcompany.%' or name = 'com.apple.loginStoreagent.plist' or name = 'com.apple.mdworker.plist' or name = 'com.apple.UIServerLogin.plist';",
        "description": "Detect RAT used by Hacking Team",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.legion.plist' or name = 'com.apple.pixel.plist';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_devilrobber_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.machook_damon.plist' or name = 'com.apple.globalupdate.plist' or name = 'com.apple.appstore.plughelper.plist' or name = 'com.apple.MailServiceAgentHelper.plist' or name = 'com.apple.systemkeychain-helper.plist' or name = 'com.apple.periodic-dd-mm-yy.plist';",
        "description": "(https://github.com/PaloAltoNetworks-BD/WireLurkerDetector)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from file where path = '/Users/Shared/UserEvent.app';",
        "description": "(http://www.intego.com/mac-security-blog/new-mac-trojan-discovered-related-to-syria/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.mdworker.plist';",
        "description": "(http://www.trendmicro.com/vinfo/us/threat-encyclopedia/malware/osx_morcut.a)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.apple.service.clipboardd.plist';",
        "description": "(https://www.fireeye.com/blog/threat-research/2014/09/forced-to-adapt-xslcmd-backdoor-now-on-os-x.html)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'realPlayerUpdate.plist';",
        "description": "(https://www.f-secure.com/weblog/archives/00002546.html)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where path like '%com.apple.AudioService.plist';",
        "description": "(http://www.intego.com/mac-security-blog/os-x-malware-tibet-variant-found/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'checkflr.plist';",
        "description": "(https://www.f-secure.com/v-descs/backdoor_osx_imuler_a.shtml)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'mac.Dockset.deman.plist';",
        "description": "(http://www.trendmicro.com/vinfo/us/threat-encyclopedia/malware/osx_dockster.a)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from file where path = '/dev/ptmx0';",
        "description": "Detect RAT used by Hacking Team",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.rxs.syslogagent.plist';",
        "description": "(http://www.symantec.com/security_response/writeup.jsp?docid=2010-081606-4034-99&tabid=2)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.codecm.uploader.plist';",
        "description": "(http://www.thesafemac.com/osxfkcodec-a-in-action/)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where name = 'com.JavaW.plist';",
        "description": "(https://www.virusbtn.com/virusbulletin/archive/2014/10/vb201410-iWorm)",
        "interval": "86400"
      },
      {
        "value": "Artifact used by this malware",
        "query": "select * from launchd where path like '%com.apple.launchport.plist';",
        "description": "(http://blog.kaspersky.com/the-mask-unveiling-the-worlds-most-sophisticated-apt-campaign/)",
        "interval": "86400"
      }
    ]
  },
  "actions": {"debug": {"sub_dir": "osx-attacks"}}
}
