# zentral/contrib/osquery/osx_package/build.tmpl/root/Library/LaunchDaemons/io.zentral.osquery.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>io.zentral.osquery.watchdog</string>
	<key>ProgramArguments</key>
	<array>
            <string>/usr/local/zentral/osquery/zentral_osquery_watchdog.py</string>
            <string>/Library/LaunchDaemons/com.facebook.osqueryd.plist</string>
            <string>/usr/local/zentral/osquery/db/</string>
            <string>/var/log/osquery/osqueryd.WARNING</string>
            <string>/usr/local/zentral/osquery/log_file_registry.plist</string>
	</array>
	<key>StartInterval</key>
        <integer>60</integer>
</dict>
</plist>
