# zentral/contrib/filebeat/osx_package/build.tmpl/root/Library/LaunchDaemons/io.zentral.filebeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>io.zentral.filebeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/zentral/bin/filebeat</string>
        <string>run</string>
        <string>-c</string>
        <string>/usr/local/zentral/filebeat/filebeat.yml</string>
        <string>-path.home</string>
        <string>/usr/local/zentral/filebeat</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
