# zentral/contrib/santa/osx_package/build.tmpl/scripts/postinstall
#!/bin/bash
## postinstall

## unload com.google.santad.plist
/bin/launchctl unload "/Library/LaunchDaemons/com.google.santad.plist"

## sleep, because …
sleep 1

## detect client device serial_number
SERIAL_NUMBER=$(ioreg -c IOPlatformExpertDevice -d 2 | awk -F\" '/IOPlatformSerialNumber/{print $(NF-1)}')
MACHINE_ID="%API_SECRET%\$SERIAL\$${SERIAL_NUMBER}"

CONFIG_PLIST="/var/db/santa/config.plist"

BLACKLIST_REGEX="%BLACKLIST_REGEX%"
WHITELIST_REGEX="%WHITELIST_REGEX%"

## config
defaults write "$CONFIG_PLIST" MachineID "$MACHINE_ID"
defaults write "$CONFIG_PLIST" ClientMode "%MODE%"
defaults write "$CONFIG_PLIST" SyncBaseURL "https://%TLS_HOSTNAME%/santa/"
defaults write "$CONFIG_PLIST" ServerAuthRootsFile "%SERVER_AUTH_ROOTS_FILE%"
if [ -z "$BLACKLIST_REGEX" ]; then
    defaults delete "$CONFIG_PLIST" BlacklistRegex
else
    defaults write "$CONFIG_PLIST" BlacklistRegex "$BLACKLIST_REGEX"
fi
if [ -z "$WHITELIST_REGEX" ]; then
    defaults delete "$CONFIG_PLIST" WhitelistRegex
else
    defaults write "$CONFIG_PLIST" WhitelistRegex "$WHITELIST_REGEX"
fi

## fix config.plist permissions
chmod 0644 /var/db/santa/config.plist

## load com.google.santad.plist
/bin/launchctl load -w "/Library/LaunchDaemons/com.google.santad.plist"

## write zentral base url
defaults write /Library/Preferences/io.zentral.plist base_url "https://%TLS_HOSTNAME%"

## force santa sync (=> preflight => machine_id => serial number)
/usr/local/bin/santactl sync

exit 0
