# zentral/contrib/santa/osx_package/build.tmpl/scripts/postinstall
#!/bin/sh
## postinstall

## detect client device serial_number
SERIAL_NUMBER=$(/usr/sbin/system_profiler SPHardwareDataType | awk '/Serial/ {print $4}')
MACHINE_ID="%API_SECRET%\$SERIAL\$${SERIAL_NUMBER}"
MACHINE_OWNER=$(/usr/bin/whoami)

## set machine_id
/usr/bin/sed -i -e "s/%MACHINE_ID%/$MACHINE_ID/" /var/db/santa/config.plist

## set machine_owner
/usr/bin/sed -i -e "s/%MACHINE_OWNER%/$MACHINE_OWNER/" /var/db/santa/config.plist

## load com.google.santad.plist
/bin/launchctl load -w "/Library/LaunchDaemons/com.google.santad.plist"  >/dev/null 2>&1

## force santa sync (=> preflight => machine_id => serial number)
/usr/local/bin/santactl sync

## write zentral base url
/usr/bin/defaults write /Library/Preferences/io.zentral.plist base_url "https://%TLS_HOSTNAME%"

exit 0
