# zentral/contrib/osquery/osx_package/build.tmpl/scripts/postinstall
#!/bin/sh
## postinstall

## load com.facebook.osqueryd
/bin/launchctl load -w "/Library/LaunchDaemons/com.facebook.osqueryd.plist"

## load io.zentral.osquery.watchdog.plist
/bin/launchctl load -w "/Library/LaunchDaemons/io.zentral.osquery.watchdog.plist"

## write zentral base url
/usr/bin/defaults write /Library/Preferences/io.zentral.plist base_url "https://%TLS_HOSTNAME%"

exit 0
