# docker/conf/probes/incident-response.json
{
  "name": "incident-response",
  "osquery": {
    "schedule": [
      {
        "interval": "3600",
        "value": "Scope for lateral movement. Potential exfiltration locations. Potential dormant backdoors.",
        "version": "1.4.5",
        "query": "select * from mounts;",
        "description": "Retrieves the current list of mounted drives in the target system."
      },
      {
        "interval": "3600",
        "value": "Detect backdoor binaries (attacker may drop a copy of /bin/sh). Find potential elevation points / vulnerabilities in the standard build.",
        "version": "1.4.5",
        "query": "select * from suid_bin;",
        "description": "Retrieves all the files in the target system that are setuid enabled."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select username, key, subkey, value from preferences p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/com.apple.loginwindow.plist';",
        "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
        "interval": "86400"
      },
      {
        "interval": "86400",
        "value": "Identify network communications that are being redirected. Example: identify if security logging has been disabled",
        "version": "1.4.5",
        "query": "select * from etc_hosts;",
        "description": "Retrieves all the entries in the target system /etc/hosts file."
      },
      {
        "platform": "darwin",
        "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
        "version": "1.4.5",
        "query": "select * from alf_explicit_auths;",
        "description": "Retrieves the list of processes with explicit authorization for the Application Layer Firewall.",
        "interval": "3600"
      },
      {
        "interval": "86400",
        "value": "Identify malware via connections to known bad IP addresses as well as odd local or remote port bindings",
        "version": "1.4.5",
        "query": "select distinct pid, family, protocol, local_address, local_port, remote_address, remote_port, path from process_open_sockets where path <> '' or remote_address <> '';",
        "description": "Retrieves all the open sockets per process in the target system."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select * from launchd;",
        "description": "Retrieves all the daemons that will run in the start of the target OSX system.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
        "version": "1.4.5",
        "query": "select * from alf_services;",
        "description": "Retrieves the services for the Application Layer Firewall in OSX.",
        "interval": "3600"
      },
      {
        "interval": "86400",
        "value": "Identify processes accessing sensitive files they shouldn't",
        "version": "1.4.5",
        "query": "select distinct pid, path from process_open_files where path not like '/private/var/folders%' and path not like '/System/Library/%' and path not in ('/dev/null', '/dev/urandom', '/dev/random');",
        "description": "Retrieves all the open files per process in the target system."
      },
      {
        "interval": "3600",
        "value": "Useful for intrusion detection and incident response. Verify assumptions of what accounts should be accessing what systems and identify machines accessed during a compromise.",
        "version": "1.4.5",
        "query": "select * from last;",
        "description": "Retrieves the list of the latest logins with PID, username and timestamp."
      },
      {
        "platform": "linux",
        "value": "Ability to compare with known good. Identify mapped regions corresponding with or containing injected code.",
        "version": "1.4.5",
        "query": "select * from process_memory_map;",
        "description": "Retrieves the memory map per process in the target Linux system.",
        "interval": "86400"
      },
      {
        "interval": "86400",
        "value": "Identify actions taken. Useful for compromised hosts.",
        "version": "1.4.5",
        "query": "select * from shell_history;",
        "description": "Retrieves the command history, per user, by parsing the shell history files."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that has a kernel extension component.",
        "version": "1.4.5",
        "query": "select * from kernel_extensions;",
        "description": "Retrieves all the information about the current kernel extensions for the target OSX system.",
        "interval": "3600"
      },
      {
        "interval": "3600",
        "value": "Identify if an attacker is using temporary, memory storage to avoid touching disk for anti-forensics purposes",
        "version": "1.4.5",
        "query": "select * from block_devices where type = 'Virtual Interface';",
        "description": "Retrieves all the ramdisk currently mounted in the target system."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select * from startup_items;",
        "description": "Retrieve all the items that will load when the target OSX system starts.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Post-priori hijack detection, detect potential sensitive information leakage.",
        "version": "1.4.7",
        "query": "select * from sandboxes;",
        "description": "Lists the application bundle that owns a sandbox label.",
        "interval": "86400"
      },
      {
        "interval": "3600",
        "value": "Identify malware that uses this persistence mechanism to launch at a given interval",
        "version": "1.4.5",
        "query": "select * from crontab;",
        "description": "Retrieves all the jobs scheduled in crontab in the target system."
      },
      {
        "platform": "darwin",
        "value": "Scope for lateral movement. Potential exfiltration locations. Potential dormant backdoors.",
        "version": "1.4.5",
        "query": "select * from nfs_shares;",
        "description": "Retrieves the current list of Network File System mounted shares.",
        "interval": "3600"
      },
      {
        "platform": "linux",
        "value": "Identify malware that has a kernel module component.",
        "version": "1.4.5",
        "query": "select * from kernel_modules;",
        "description": "Retrieves all the information for the current kernel modules in the target Linux system.",
        "interval": "3600"
      },
      {
        "interval": "3600",
        "value": "Detect if a listening port iis not mapped to a known process. Find backdoors.",
        "version": "1.4.5",
        "query": "select * from listening_ports;",
        "description": "Retrieves all the listening ports in the target system."
      },
      {
        "platform": "darwin",
        "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
        "version": "1.4.5",
        "query": "select * from alf_exceptions;",
        "description": "Retrieves the exceptions for the Application Layer Firewall in OSX.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Post-priori hijack detection, detect potential sensitive information leakage.",
        "version": "1.4.7",
        "query": "select * from app_schemes;",
        "description": "Retreives the list of application scheme/protocol-based IPC handlers.",
        "interval": "86400"
      },
      {
        "interval": "3600",
        "value": "Detect trojanized system accounts,ie, password change for 'bin'. New attacker owned accounts.",
        "version": "1.4.5",
        "query": "select * from passwd_changes;",
        "description": "Retrieves all the password changes in the target system."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select key, subkey, value from preferences where path = '/Library/Preferences/loginwindow.plist';",
        "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
        "interval": "86400"
      },
      {
        "interval": "86400",
        "value": "Determine if a third party device was attached to the system, include in scope for identification/mitigation/future detection.",
        "version": "1.4.5",
        "query": "select distinct * from hardware_events where path <> '' or model <> '';",
        "description": "Retrieves all the hardware related events in the target OSX system."
      },
      {
        "interval": "86400",
        "value": "Identifies a system potentially vulnerable to disk cloning.",
        "version": "1.4.5",
        "query": "select * from disk_encryption;",
        "description": "Retrieves the current disk encryption status for the target system."
      },
      {
        "interval": "3600",
        "value": "Useful for intrusion detection and incident response. Verify assumptions of what accounts should be accessing what systems and identify machines accessed during a compromise.",
        "version": "1.4.5",
        "query": "select liu.*, p.name, p.cmdline, p.cwd, p.root from logged_in_users liu, processes p where liu.pid = p.pid;",
        "description": "Retrieves the list of all the currently logged in users in the target system."
      },
      {
        "platform": "linux",
        "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
        "version": "1.4.5",
        "query": "select * from iptables;",
        "description": "Retrieves the current filters and chains per filter in the target system.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Identifies connections to rogue access points.",
        "version": "1.4.5",
        "query": "select key, subkey, value from preferences where path = '/Library/Preferences/SystemConfiguration/com.apple.airport.preferences.plist';",
        "description": "Retrieves all the remembered wireless network that the target machine has connected to.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Identify recently accessed items. Useful for compromised hosts.",
        "version": "1.4.5",
        "query": "select username, key, value from preferences p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/com.apple.recentitems.plist';",
        "description": "Retrieves the list of recent items opened in OSX by parsing the plist per user.",
        "interval": "86400"
      },
      {
        "interval": "3600",
        "value": "Determine if MITM in progress.",
        "version": "1.4.5",
        "query": "select * from arp_cache;",
        "description": "Retrieves the ARP cache values in the target system."
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select username, key, subkey, value from preferences p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/loginwindow.plist';",
        "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
        "interval": "86400"
      },
      {
        "platform": "darwin",
        "value": "Determine what access the user had, identify accounts to be remediated.",
        "version": "1.4.5",
        "query": "select * from keychain_items;",
        "description": "Retrieves all the items contained in the keychain in the target OSX system.",
        "interval": "3600"
      },
      {
        "interval": "3600",
        "value": "Identify if a machine is being used as relay",
        "version": "1.4.5",
        "query": "select * from system_controls where name = 'net.inet.ip.forwarding';",
        "description": "Retrieves the current status of IP forwarding."
      },
      {
        "platform": "darwin",
        "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
        "version": "1.4.5",
        "query": "select * from alf;",
        "description": "Retrieves the configuration values for the Application Layer Firewall for OSX.",
        "interval": "3600"
      },
      {
        "interval": "86400",
        "value": "Insight into the process data: Where was it started from, was it preloaded...",
        "version": "1.4.5",
        "query": "select * from process_envs;",
        "description": "Retrieves all the environment variables per process in the target system."
      },
      {
        "platform": "darwin",
        "value": "Identify malware, adware, or vulnerable packages that are installed as an application.",
        "version": "1.4.5",
        "query": "select * from apps;",
        "description": "Retrieves all the currently installed applications in the target OSX system.",
        "interval": "3600"
      },
      {
        "platform": "darwin",
        "value": "Identify malware that uses this persistence mechanism to launch at system boot",
        "version": "1.4.5",
        "query": "select key, subkey, value from preferences where path = '/Library/Preferences/com.apple.loginwindow.plist';",
        "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
        "interval": "86400"
      }
    ]
  },
  "actions": {"debug": {"sub_dir": "incident-response"}}
}
