# zentral/contrib/monolith/osx_package/build.tmpl/root/Library/LaunchAgents/io.zentral.monolith.depnotify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.zentral.monolith.depnotify</string>
	<key>ProgramArguments</key>
	<array>
            <string>/Applications/Utilities/DEPNotify.app/Contents/MacOS/DEPNotify</string>
		<string>-munki</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
